# roles/ansible-anaconda/templates/etc/profile.d/anaconda.sh.j2
export PATH={{anaconda_parent_dir}}/anaconda{{anaconda_python_ver}}/bin:$PATH
export PYTHON={{anaconda_parent_dir}}/anaconda{{anaconda_python_ver}}/bin/python
